# templates/systemd.service.erb
[Service]

# Maximum number of open files (for ulimit -n)
LimitNOFILE=<%= @ulimit %>
